# hawks.conf
[program:hawks]
directory=/home/username/encratica_gmail_com/dev/hawks_dash
command=/home/username/.venvs/hawks/bin/python app.py
